# fuel-sync-watcher.service
[Unit]
Description=Fuel Centre Firestore Sync Watcher
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/alpha/fuel-dev
EnvironmentFile=/home/alpha/.env/fuel.env
ExecStart=/usr/bin/node scripts/firestore-sync.mjs watch
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
